# supervisor.ini
[program:huobi_server_spot]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/huobi_server_spot.php start
stdout_logfile=/var/www/app/websocket/log_huobi_server_spot.log

[program:huobi_client_spot]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/huobi_client_spot.php start
stdout_logfile=/var/www/app/websocket/log_huobi_client_spot.log


[program:huobi_server_future]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/huobi_server_future.php start
stdout_logfile=/var/www/app/websocket/log_huobi_server_future.log

[program:huobi_client_future]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/huobi_client_future.php start
stdout_logfile=/var/www/app/websocket/log_huobi_client_future.log

[program:okex_server]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/okex_server.php start
stdout_logfile=/var/www/app/websocket/log_okex_server.log

[program:okex_client]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/okex_client.php start
stdout_logfile=/var/www/app/websocket/log_okex_client.log

[program:bitmex_server]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/bitmex_server.php start
stdout_logfile=/var/www/app/websocket/log_bitmex_server.log

[program:bitmex_client]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/bitmex_client.php start
stdout_logfile=/var/www/app/websocket/log_bitmex_client.log

[program:binance_server_spot]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/binance_server_spot.php start
stdout_logfile=/var/www/app/websocket/log_binance_server_spot.log

[program:binance_client_spot]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/websocket/binance_client_spot.php start
stdout_logfile=/var/www/app/websocket/log_binance_client_spot.log


